# linux/omarchy/qbittorrent-nox.service
[Unit]
    Description=qBittorrent-nox service
    After=network.target

    [Service]
    Type=forking
    User=taxin
    ExecStart=/usr/bin/qbittorrent-nox -d --webui-port=8080
    Restart=on-failure

    [Install]
    WantedBy=multi-user.target
    ```

3.  **Enable and start the service**:
    
```bash
    sudo systemctl daemon-reload
    sudo systemctl enable --now qbittorrent-nox
    ```

### 3. Verify the Status
You can now check if it's running correctly using the alias you likely have or by running:
```bash
systemctl status qbittorrent-nox
